# files/default/rsync-files/rsync-files.cron
# Sync files
*/5 * * * * root /etc/rsync-files/rsync-files.sh

